# mt7620/ImageBuilder/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
opkg install /etc/luci-theme-darkmatter_git-17.005.53085-cca792b-1_all.ipk
rm -f /etc/luci-theme-darkmatter_git-17.005.53085-cca792b-1_all.ipk
/etc/init.d/haproxy disable
/etc/init.d/haproxy stop
/etc/init.d/pdnsd disable
/etc/init.d/pdnsd stop
echo exit 0 >/etc/rc.local

exit 0
